# scripts/com.pallotta.llama-background.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pallotta.llama-background</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/pallotta/Personal AI Assistant/scripts/start_llama_background_monitored.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/llama-background-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/llama-background-launchd-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/pallotta/Personal AI Assistant</string>
</dict>
</plist>

